# configs/launchd/com.quant-system.console.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
     <key>Label</key>
     <string>com.quant-system.console</string>

      <key>ProgramArguments</key>
     <array>
          <string>/Users/ziyuan/Code/quant-system/.venv/bin/uv</string>
          <string>run</string>
          <string>quant</string>
          <string>web</string>
          <string>serve</string>
          <string>--host</string>
          <string>127.0.0.1</string>
          <string>--port</string>
          <string>8000</string>
     </array>

      <key>WorkingDirectory</key>
     <string>/Users/ziyuan/Code/quant-system</string>

      <key>EnvironmentVariables</key>
     <dict>
          <key>QUANT_CONSOLE_API_KEY</key>
          <string>CHANGE_ME_TO_A_RANDOM_VALUE</string>
          <key>PATH</key>
          <string>/Users/ziyuan/Code/quant-system/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
     </dict>

      <key>StandardOutPath</key>
     <string>/Users/ziyuan/Code/quant-system/logs/console.log</string>

      <key>StandardErrorPath</key>
     <string>/Users/ziyuan/Code/quant-system/logs/console-error.log</string>

      <key>KeepAlive</key>
     <true/>

      <key>RunAtLoad</key>
     <true/>

      <key>Disabled</key>
     <true/>
</dict>
</plist>
