# configs/launchd/com.quant-system.alpaca-paper-refresh.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quant-system.alpaca-paper-refresh</string>

  <key>Disabled</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/absolute/path/to/quant-system</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/absolute/path/to/quant-system/scripts/run_alpaca_paper_refresh.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/absolute/path/to/quant-system/logs/launchd-alpaca-paper-refresh.out.log</string>

  <key>StandardErrorPath</key>
  <string>/absolute/path/to/quant-system/logs/launchd-alpaca-paper-refresh.err.log</string>
</dict>
</plist>
